# .config/systemd/user/mopidy-user.service
[Unit]
Description=Mopidy music server

[Service]
Type=simple
ExecStart=/usr/bin/mopidy &
Restart=always

[Install]
WantedBy=default.target
