# .local/share/dotfiles/home-backup/.config/systemd/user/polybar.service
[Unit]
Description=Polybar

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=%h/.config/polybar/launch.sh
ExecStop=killall -q polybar

[Install]
WantedBy=multi-user.target
